# docker/crontab
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin
*/30 * * * * root /app/scripts/home_fetch.sh >> /var/log/bicing-fetch.log 2>&1
