# docker/crontab
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin
*/10 * * * * root cd /app && /usr/local/bin/python scripts/ingest.py >> /var/log/bicing-ingest.log 2>&1 && /usr/local/bin/python scripts/export.py >> /var/log/bicing-export.log 2>&1 && /app/scripts/deploy.sh >> /var/log/bicing-deploy.log 2>&1
